# deploy/relay.service
[Unit]
Description=Telegram Relay Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-relay
ExecStart=/opt/telegram-relay/relay
Restart=always
RestartSec=3
User=relay
Group=relay
Environment=LISTEN_ADDR=:8780
Environment=SECURITY_LEVEL=medium
Environment=AUTH_TOKEN=replace-with-strong-token
Environment=HMAC_SECRET=replace-with-hmac-secret
Environment=TIMESTAMP_SKEW_SEC=300
Environment=RATE_LIMIT_PER_SECOND=20
Environment=RATE_LIMIT_BURST=40
Environment=TELEGRAM_BOT_TOKEN=replace-with-bot-token
Environment=TELEGRAM_CHAT_ID=-1001234567890
Environment=TELEGRAM_PARSE_MODE=HTML
Environment=TELEGRAM_API_BASE_URL=https://api.telegram.org
Environment=TELEGRAM_TIMEOUT_SEC=5
Environment=RETRY_MAX_ATTEMPTS=3
Environment=RETRY_INITIAL_BACKOFF_MS=300
Environment=RETRY_MAX_BACKOFF_MS=3000

[Install]
WantedBy=multi-user.target
